# files/sslbackup-full.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
55 23 * * 7 root /usr/bin/sslbackup.sh full
