# files/sslbackup-diff.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
55 23 * * 1-6 root /usr/bin/sslbackup.sh diff
